# ansible/roles/syncthing/files/syncthing@syncthing.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/syncthing serve --no-browser --no-restart --logflags=0 --home=/mnt/storage/appdata/syncthing
