# ansible/roles/traefik/files/traefik.service
[Unit]
Description=Traefik service

[Service]
User=traefik
Group=traefik
ExecStart=traefik

[Install]
WantedBy=multi-user.target
